# cc-proxy.service
[Unit]
Description=Claude Code Proxy (mitmproxy)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/danny/code/cc-proxy
ExecStart=/home/danny/code/cc-proxy/cc-proxy
Restart=on-failure
RestartSec=5

# Logging
StandardOutput=append:%h/.local/log/cc-proxy/cc-proxy.log
StandardError=append:%h/.local/log/cc-proxy/cc-proxy.log

[Install]
WantedBy=default.target
